# modules/systemd-hardening/tests/fixtures/systemd-units/onionleak-collector.service
# Fixture only — minimal onionleak-collector unit for the bats suite.
[Unit]
Description=onionleak collector (fixture)

[Service]
Type=simple
ExecStart=/usr/bin/onionleak-collector
Restart=on-failure

[Install]
WantedBy=multi-user.target
